# infra/openclaw/alcomatcher-openclaw-digest.timer
[Unit]
Description=AlcoMatcher OpenClaw KPI Digest Timer

[Timer]
OnCalendar=*-*-* 08:00:00
OnCalendar=*-*-* 18:00:00
Persistent=true
Unit=alcomatcher-openclaw-digest.service

[Install]
WantedBy=timers.target
