# infra/openclaw/alcomatcher-openclaw-monitor.service
[Unit]
Description=AlcoMatcher OpenClaw Ops Monitor (tiered alerts)
After=docker.service network.target
Requires=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/alcomatcher
Environment=TELEGRAM_TARGET=593360085
Environment=API_BASE=http://127.0.0.1:3000
Environment=ALERT_STATE_DIR=/var/tmp/alcomatcher-openclaw/state
Environment=POLL_SECONDS=60
Environment=WINDOW_SECONDS=900
Environment=HEALTH_FAIL_THRESHOLD=3
Environment=HTTP_5XX_THRESHOLD=20
Environment=SCANNER_FAIL_THRESHOLD=5
Environment=SYNC_FAILED_DELTA_THRESHOLD=5
Environment=TELEMETRY_COMPLETE_MIN_PCT=90
Environment=SLO_P50_LIMIT_MS=5000
Environment=SLO_P95_LIMIT_MS=8000
ExecStart=/bin/bash /opt/alcomatcher/infra/openclaw/monitor-alerts.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
